# supporting_files/supervisord-syncthing.conf
[program:syncthing]
command=/start-syncthing.sh
numprocs=1
autostart=true
autorestart=true
